# config/service/pipewire-pulse/run
#!/bin/sh
exec 2>&1

# Check wireplumber
sv check "$HOME/.config/service/wireplumber" || exit 1

exec pipewire-pulse
